# etc/supervisor/conf.d/sploshify_player.conf
[program:sploshify_player]
command=/home/pi/sploshify/sploshify_player.py
autostart=true
autorestart=true
stderr_logfile=/var/log/sploshify.log
stdout_logfile=/var/log/sploshify.log
